# scripts/rules/60-device-trigger.rules
ACTION=="add", KERNEL=="mmcblk0p2", RUN+="/usr/bin/udevadm trigger --type=devices"
