# scripts/rules/60-device-trigger.rules
ACTION=="add", KERNEL=="mmcblk[0-9]p[0-9]", RUN+="/usr/bin/udevadm trigger --type=devices"
