# scripts/rules/99-adb-restart.rules
ACTION=="change", SUBSYSTEM=="udc", RUN+="/usr/bin/restart_adbd.sh"
