# scripts/units/serial-getty@.service.d/deps-booting-done.conf
# Launching login shell after booting-done service

[Unit]
After=booting-done.service
Requires=booting-done.service
